# docker/s6-overlay/s6-rc.d/svc-firewall-refresh/run
#!/command/with-contenv /bin/bash
exec /etc/s6-overlay/scripts/firewall-refresh.sh
