# docker/s6-overlay/s6-rc.d/svc-bootstrap/up
#!/command/with-contenv /bin/bash
exec s6-setuidgid developer /etc/s6-overlay/scripts/bootstrap.sh
