# docker/s6-overlay/s6-rc.d/init-workspace-check/up
#!/command/with-contenv /bin/bash
exec /etc/s6-overlay/scripts/init-workspace-check.sh
